# myDoorbell.cron
* * * * * /usr/local/bin/python /home/myDoorbell/bin/myDoorbellConfig.py 1>> /tmp/myDoorbellConfig.log
